# docker/supervisor/conf.d/laravel-schedule.conf
[program:laravel-schedule]
command=/bin/sh -c "while true; do php /var/www/app/artisan schedule:run --verbose --no-interaction; sleep 60; done"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/schedule.log
priority=3
